# puppet/environments/production/modules/tdp_app/templates/tdp-recruitment.service.erb
[Unit]
Description=Java Service

[Service]
User=root
# The configuration file application.properties should be here:
WorkingDirectory=/etc/tdp-recruitment/app_config.yml
ExecStart=/usr/bin/java -Xmx256m -jar /opt/tdp-recruitment/tdp-recruitment-1.0-SNAPSHOT.jar server
SuccessExitStatus=143

[Install]
WantedBy=multi-user.target
